# root/etc/services.d/fail2ban/run
#!/usr/bin/with-contenv bash

fail2ban-client -x start
