# root/etc/services.d/fail2ban/run
#!/usr/bin/with-contenv bash
# shellcheck disable=SC1008

exec fail2ban-client -x -f start
